# tovest-bot.service
[Unit]
Description=Tovest Tracking Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/tovest_tracking_bot
EnvironmentFile=/home/ubuntu/tovest_tracking_bot/.env
ExecStart=/home/ubuntu/tovest_tracking_bot/venv/bin/python bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
